# foobar.target
[Unit]
Description=foo and bar

[Install]
WantedBy=default.target
